# WatchMe project/pro.v2/infra/systemd/watchme-v2-api.service
[Unit]
Description=WatchMe V2 API server
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=watchme
Group=watchme
WorkingDirectory=/opt/watchme-v2
EnvironmentFile=-/etc/watchme-v2/api.env
ExecStart=/usr/bin/node apps/api/src/index.js
Restart=always
RestartSec=5
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target
