# WatchMe project/pro.v2/infra/systemd/watchme-v2-pg-backup.service
[Unit]
Description=WatchMe V2 Postgres logical backup (pg_dump)
Wants=network-online.target

[Service]
Type=oneshot
User=watchme
Group=watchme
WorkingDirectory=/opt/watchme-v2
EnvironmentFile=-/etc/watchme-v2/backup.env
ExecStart=/usr/bin/node scripts/pg-backup.js --outDir=/var/backups/watchme-v2/pg

[Install]
WantedBy=multi-user.target
