# WatchMe project/pro.v2/infra/systemd/watchme-v2-worker.service
[Unit]
Description=WatchMe V2 worker (queue runner)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=watchme
Group=watchme
WorkingDirectory=/opt/watchme-v2
EnvironmentFile=-/etc/watchme-v2/worker.env
ExecStart=/usr/bin/node apps/worker/src/index.js
Restart=always
RestartSec=5
TimeoutStopSec=60

[Install]
WantedBy=multi-user.target
